# scripts/openclaw-gateway-launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.gateway</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/zhonglian/.nvm/versions/node/v22.22.1/bin/node</string>
    <string>/Users/zhonglian/openclaw/openclaw.mjs</string>
    <string>gateway</string>
    <string>run</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/zhonglian/openclaw</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OPENCLAW_BAIDU_PROJECT_ROOT</key>
    <string>/Volumes/ragflow/hotmaxx/hotmaxxflag</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/zhonglian/.openclaw/logs/gateway.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/zhonglian/.openclaw/logs/gateway.err.log</string>
</dict>
</plist>
